# systemd/govscout-sync.timer
[Unit]
Description=Run GovScout Sync daily at 2am

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
